# Server/docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:redis]
command=redis-server --save "" --appendonly no --loglevel warning
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:server]
command=sh -c "until redis-cli ping; do sleep 0.5; done && cd /app && python -m alembic upgrade head && uvicorn src.main:app --host 0.0.0.0 --port %(ENV_PORT)s"
environment=PORT="%(ENV_PORT)s"
autostart=true
autorestart=true
startretries=5
startsecs=5
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
